# certbot/root/etc/services.d/certbot/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

exec snooze -D/7 /scripts/certbot_svc
