# taskserver/root/etc/services.d/taskd/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

exec taskd server --data "$TASKDDATA"
