# taskserver/root/etc/services.d/taskd/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

# exec taskd server --debug --debug.tls=2 # --data "$TASKDDATA"
exec s6-setuidgid taskd taskd server --debug --debug.tls=2
